# openprocurement/auction/templates/systemd.timer
[Unit]
Description=Auction Timer for Auction {{ description }}

[Timer]
OnCalendar={{ timestamp }}
 
[Install]
WantedBy=timers.target